# package/base-files/files/etc/crontabs/root
*/2 * * * * isfound=$(ps | grep "/usr/sbin/ccnu" | grep -v "grep"); if [ -z "$isfound" ]; then /etc/init.d/ccnu restart; fi