# package/ralink/applications/hwnat/files/hwnat
#!/bin/sh /etc/rc.common

START=91
STOP=91

USE_PROCD=1

start_service() {
	insmod /lib/modules/ralink/hw_nat.ko
	mknod /dev/hwnat0 c 220 0
}

stop_service() {
	rmmod hw_nat.ko
}


