# package/external/chengr28/pcap-dnsproxy/files/pcap-dnsproxy.init
#!/bin/sh /etc/rc.common

START=95

SERVICE_USE_PID=1
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

#TODO: startup params
# CONFIG=

start() {
	service_start /usr/sbin/Pcap_DNSProxy
}

stop() {
	service_stop /usr/sbin/Pcap_DNSProxy
}
