# package/external/chengr28/pcap-dnsproxy/files/pcap-dnsproxy.init
#!/bin/sh /etc/rc.common

START=95

SERVICE_USE_PID=1
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

# MUST set Config.conf before 'enable' this service

start() {
	service_start /usr/sbin/Pcap_DNSProxy -c /etc/pcap-dnsproxy
}

stop() {
	# FIXME: something wrong with 'service_stop'
	killall -SIGTERM Pcap_DNSProxy
}
