# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
# enable tcp fast open
echo 3 > /proc/sys/net/ipv4/tcp_fastopen
exit 0
