# target/linux/mvebu/base-files/etc/crontabs/root
*/5 * * * * /sbin/fan_ctrl.sh >/dev/null
